# .global.d/Library/LaunchDaemons/dev.artero.sysctl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.artero.sysctl</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/sbin/sysctl</string>
      <string>kern.maxfiles=1048576</string>
      <string>kern.maxfilesperproc=65536</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
