# .systems/darwin/LaunchAgents/Library/LaunchAgents/dev.artero.environment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.artero.environment</string>
    <key>LimitLoadToSessionType</key>
    <array>
      <string>Aqua</string>
      <string>StandardIO</string>
      <string>Background</string>
    </array>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>/Users/william/Library/Scripts/dotenv.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>WatchPaths</key>
    <array>
      <string>/Users/william/.env</string>
      <string>/Users/william/.env_remove</string>
      <string>/Users/william/.env_secrets</string>
      <string>/Users/william/.env_personal</string>
      <string>/Users/william/.env_personal_secrets</string>
      <string>/Users/william/.env_work</string>
      <string>/Users/william/.env_work_secrets</string>
      <string>/Users/william/.env_remove_work_secrets</string>
    </array>
  </dict>
</plist>
