# .hostnames/C02DQ36NMD6P/tinyproxy/Library/LaunchAgents/dev.artero.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.artero.proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/tinyproxy</string>
    <string>-d</string>
    <string>-c</string>
    <string>/Users/william/.config/tinyproxy/tinyproxy.conf</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/william/Library/Logs/tinyproxy-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/william/Library/Logs/tinyproxy-stderr.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
